# ops/systemd/oh-my-github-runner.service
[Unit]
Description=oh-my-github-runner daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/oh-my-github-runner
EnvironmentFile=/etc/oh-my-github-runner/runner.env
ExecStart=/usr/bin/node /home/ubuntu/oh-my-github-runner/dist/src/index.js
Restart=always
RestartSec=5
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
